# print_agent/print_agent.service
[Unit]
Description=Odoo Print Agent
After=network.target

[Service]
Type=simple
User=gristmill
WorkingDirectory=/home/gristmill/Fulfillment_Automation/print_agent
ExecStart=/home/gristmill/Fulfillment_Automation/print_agent/venv/bin/python3 main.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target


